# systemd-units/dedicated-server-stop.timer
[Unit]
Description=timer for stopping dedicated server

[Timer]
OnActiveSec=15min
