# systemd-units/mcserv.socket
[Unit]
Description=Socket for mcserv.service
Conflicts=mcserv.service

[Socket]
ListenStream=25565

[Install]
WantedBy=network-online.target
