# systemd-units/mcserv.socket
[Unit]
Description=Socket for Starting Minecraft Server

[Socket]
ListenStream=25565

[Install]
WantedBy=sockets.target
