# systemd-units/mcserv-stop.service
[Unit]
Description=Stop Minecraft Server

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl stop mcserv.service
