# systemd-units/mcserv-stop.timer
[Unit]
Description=Timer for Stopping Minecraft server

[Timer]
OnActiveSec=15min

[Install]
WantedBy=multi-user.target
