# systemd-units/mcserv.service
[Unit]
Description=Minecraft Server
After=network-online.target

Wants=mcserv-stoptimerctrl.service
OnFailure=mcserv.socket

[Service]
WorkingDirectory=/opt/mcserv
ExecStart=/usr/local/bin/mcserv-start.sh
ExecStopPost=systemctl start mcserv.socket

[Install]
WantedBy=mutli-user.target
