# systemd-units/mcserv.service
[Unit]
Description=Minecraft Server
After=network-online.target
Wants=mcserv-stoptimerctrl.service
Conflicts=mcserv.socket

OnFailure=mcserv.socket

[Service]
User=minecraft
Group=minecraft
PassEnvironment=MINECRAFT_SERVER_JAR MINECRAFT_JVM_ARGS MINECRAFT_JAR_ARGS
WorkingDirectory=/opt/mcserv
SuccessExitStatus=143
Restart=on-failure

ExecStart=java $MINECRAFT_JVM_ARGS -jar ${MINECRAFT_SERVER_JAR} $MINECRAFT_JAR_ARGS
ExecStopPost=/bin/false

[Install]
WantedBy=mutli-user.target
