# systemd-units/mcserv-stoptimerctrl.service
[Unit]
Description=player-count-monitor for minecraft server
BindsTo=mcserv.service

[Service]
WorkingDirectory=/opt/mcserv
ExecStart=/opt/mcserv/scripts/shutdown-timerctrl.sh
