# systemd-units/mcserv-stop.service
[Unit]
Description=service to stop mcserv

[Service]
ExecStart=/usr/bin/systemctl stop mcserv.service
