# systemd-units/mcserv-stop.service
[Unit]
Description=service to stop mcserv
Conflicts=mcserv.service

[Service]
Type=oneshot
ExecStart=/bin/echo "stopping mcserv"
